# systemd/beacon.service
[Unit]
Description=Beacon Agent - Lightweight deployment and reporting for IoT
After=network.target

[Service]
Environment="GIT_SSH_COMMAND=ssh -i /etc/beacon/id_ed25519 -o IdentitiesOnly=yes -o StrictHostKeyChecking=no"
Environment="BEACON_GIT_TOKEN=ghp_xxxxxxxx"
Environment="BEACON_REPO_URL=https://${BEACON_GIT_TOKEN}@github.com/yourusername/yourrepo.git"
Type=simple
ExecStart=/usr/local/bin/beacon
WorkingDirectory=/opt/beacon
Restart=always
RestartSec=5
User=pi

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target