# systemd/beacon@.service
# Example systemd service template for Beacon project agents.
#
# DO NOT install this file directly — beacon bootstrap generates per-project
# service files with the correct absolute paths for WorkingDirectory,
# EnvironmentFile, and HOME.
#
# This file is a reference only. To set up a project:
#   beacon bootstrap
#
# Or to create a master service:
#   beacon cloud login --api-key usr_xxx
#   beacon master  (will daemonize automatically)

[Unit]
Description=Beacon Agent for %i - Lightweight deployment and monitoring
After=network.target

[Service]
# These paths must be absolute — edit to match your setup.
# EnvironmentFile=/home/youruser/.beacon/config/projects/%i/env
Type=simple
ExecStart=/usr/local/bin/beacon deploy
# WorkingDirectory=/home/youruser/projects/%i
Restart=always
RestartSec=5

# Critical: HOME must be set for beacon to find ~/.beacon/config.yaml
# Environment=HOME=/home/youruser

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
